# sv/udev/run
#!/bin/sh
exec /sbin/udevd --children-max=2
